# deploy/homesec.service
[Unit]
Description=HomeSec control panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/homesec/backend
EnvironmentFile=/opt/homesec/backend/.env
ExecStart=/opt/homesec/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
