# deploy/homesec-update.timer
[Unit]
Description=Run HomeSec pull-deploy every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Unit=homesec-update.service

[Install]
WantedBy=timers.target
