# deploy/homesec-update.service
[Unit]
Description=HomeSec pull-deploy from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/homesec/deploy/update.sh
